# conf/systemd/fai-install.target

[Unit]
Description=Install system using FAI
Requires=sysinit.target fai-install.service
Conflicts=rescue.service rescue.target
AllowIsolate=yes
