# conf/systemd/fai.target
[Unit]
Description=Fully Automatic Installation system
Requires=sysinit.target network-online.target fai.service
AllowIsolate=no
