# conf/systemd/fai.target
[Unit]
Description=Install system using FAI
Requires=network-online.target fai-install.service
AllowIsolate=no
