# debian/fai-nfsroot.postinst
#!/bin/sh -e

case "$1" in
    configure)
	update-alternatives --install /usr/sbin/policy-rc.d policy-rc.d \
	    /usr/bin/policy-rc.d.fai 100
    update-rc.d fai-abort start 01 2 >/dev/null
esac

#DEBHELPER#

exit 0
