# debian/fai-nfsroot.postrm
#!/bin/sh -e

case "$1" in
    remove)
	rm /etc/init.d/rcS
	dpkg-divert --package fai-nfsroot --remove --rename /etc/init.d/rcS
	;;

    purge)
	update-rc.d fai-abort remove > /dev/null
	;;
esac

#DEBHELPER#

exit 0
