# debian/fai-nfsroot.postrm
#!/bin/sh -e

case "$1" in
    purge)
        update-rc.d fai-abort remove > /dev/null
        ;;
esac

#DEBHELPER#

exit 0
