# 3.4/conf/upstart-fai.conf
# fai - start the fully automatic installation

description	"start fully automatic installation"

start on local-filesystems

console output

task
exec /usr/sbin/fai
